# agenten/syslog-ng/oracle/syslog-ng2oracle.init-sles
#!/bin/sh
# Start/stop the syslog2oracle daemon.
#
# chkconfig: 345 99 01
# description: Icinga network monitor
### BEGIN INIT INFO
# Provides:          syslog-ng2oracle
# Required-Start:
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: EventDB syslog-ng2oracle Initscript
# Description:       EventDB Syslog to Oracle
### END INIT INFO

DAEMON=/usr/local/icinga/contrib/syslog-ng2oracle.pl
PID=./syslog-ng2oracle.pid
#FIFO="/var/run/syslog-ng.pipe"
FIFO=$(cat $DAEMON |egrep -i '^my \$fifo' |cut -d"=" -f2 |sed 's/"//g; s/;//g')
test -p $FIFO || mkfifo $FIFO

test -f /$DAEMON || exit 0

case "$1" in
start) echo -n "Starting syslog2oracle daemon"
	startproc -p $PID $DAEMON
	echo "."
	;;
stop) echo -n "Stoping syslog2oracle daemon"
	$DAEMON --stop --quiet --pidfile $PID
	killall syslog-ng2oracle.pl
	echo "."
	;;
*) echo "Usage: /etc/init.d/syslogd2oracle start|stop"
	exit 1
	;;
esac
exit 0

